# launchd/com.dasauto.check-listings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.dasauto.check-listings</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/anaconda3/bin/python3</string>
		<string>/Users/michael/projects/Das_Auto/scripts/check_new_listings.py</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/michael/projects/Das_Auto</string>

	<key>StartInterval</key>
	<integer>1800</integer>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/tmp/das-auto-check.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/das-auto-check-error.log</string>
</dict>
</plist>
